# playbooks/templates/10-wifi-disable-powermanagement.rules
KERNEL=="wl*", ACTION=="add", RUN+="/sbin/iwconfig {{ interface_wlan }} power off"
